# flask_app/recasketch.conf
description "Gunicorn application server running recasketch"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data

env PATH=/home/ubuntu/miniconda2/envs/recasketch/bin
chdir /home/ubuntu/recasketch/flask_app
exec  gunicorn --workers 3 --bind unix:recasketch.sock -m 007 wsgi:app

